# packages/gost/contents/usr/lib/systemd/system/gost.service
[Unit]
Description=gost service
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
User=gost
Group=gost
Restart=on-failure
RestartSec=30
ExecStart=/usr/bin/gost -C /etc/gost.yml
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
