# packages/prometheus-community/PushProx/contents/usr/lib/systemd/system/pushprox-client.service
[Unit]
Description=PushProx client
Documentation=https://github.com/prometheus-community/PushProx
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
Restart=on-failure
User=prometheus
Group=prometheus
EnvironmentFile=/etc/default/pushprox-client
ExecStart=/usr/bin/pushprox-client $PUSHPROX_CLIENT_OPTS

[Install]
WantedBy=multi-user.target
