# packages/prometheus-community/PushProx/contents/usr/lib/systemd/system/pushprox-proxy.service
[Unit]
Description=PushProx proxy
Documentation=https://github.com/prometheus-community/PushProx
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
Restart=on-failure
User=prometheus
Group=prometheus
EnvironmentFile=-/etc/default/pushprox-proxy
ExecStart=/usr/bin/pushprox-proxy $PUSHPROX_PROXY_OPTS

[Install]
WantedBy=multi-user.target
