# packages/ory/oathkeeper/contents/usr/lib/systemd/system/oathkeeper@.service
[Unit]
Description=%I oathkeeper daemon service
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
User=ory
Group=ory
Restart=on-failure
RestartSec=30
WorkingDirectory=/etc/ory/oathkeeper/%i
ExecStart=/usr/bin/oathkeeper serve --config config.yml

[Install]
WantedBy=multi-user.target
