# packages/oauth2-proxy/contents/usr/lib/systemd/system/oauth2-proxy@.service
[Unit]
Description=%I oauth2-proxy daemon service
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
User=oauth2-proxy
Group=oauth2-proxy
EnvironmentFile=-/etc/oauth2-proxy/%i.env
Restart=on-failure
RestartSec=30
WorkingDirectory=/etc/oauth2-proxy
ExecStart=/usr/bin/oauth2-proxy --config=/etc/oauth2-proxy/%i.toml

[Install]
WantedBy=multi-user.target
