# packages/ory/hydra/contents/usr/lib/systemd/system/hydra@.service
[Unit]
Description=%I hydra daemon service
After=network.target network-online.target nss-lookup.target basic.target
Wants=network-online.target nss-lookup.target

[Service]
User=ory
Group=ory
Restart=on-failure
RestartSec=30
WorkingDirectory=/etc/ory/hydra
EnvironmentFile=/etc/default/hydra
ExecStart=/usr/bin/hydra serve all --config /etc/ory/hydra/%i.yml $HYDRA_OPTS

[Install]
WantedBy=multi-user.target
